# agent/systemd/argus-agent.user.service
[Unit]
Description=ArgusMonitor host agent (user)
After=default.target network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/argus-agent/argus-agent.env
ExecStart=%h/.local/bin/argus-agent
Restart=always
RestartSec=5
WorkingDirectory=%h/.local/share/argus-agent

[Install]
WantedBy=default.target
