# agent/systemd/vordr-agent.service
[Unit]
Description=Vordr host agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/vordr-agent/vordr-agent.env
ExecStart=/usr/local/bin/vordr-agent
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/vordr-agent
StateDirectory=vordr-agent
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/vordr-agent

[Install]
WantedBy=multi-user.target
